# docker/php/supervisor.app.ini
[supervisord]
logfile=/proc/1/fd/1
logfile_maxbytes = 0
loglevel = warn
user=root

[supervisorctl]
username = dummy
password = dummy

[program:app-queue]
process_name=%(program_name)s_%(process_num)02d
command=bin/console messenger:consume async -vv --time-limit=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
stopwaitsecs=60
stdout_logfile=/dev/null
stdout_logfile_maxbytes = 0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes = 0
environment=MONOLOG_LOG_STREAM="/app/var/log/messenger.log"

[program:app-cron]
process_name=%(program_name)s_%(process_num)02d
command=bin/console messenger:consume scheduler_default -vv
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
stopwaitsecs=60
stdout_logfile=/dev/null
stdout_logfile_maxbytes = 0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes = 0
environment=MONOLOG_LOG_STREAM="/app/var/log/cron.log"
